# tests/sqlite-demo/Ubuntu_20.04/sqlite-demo/DEBIAN/prerm
#!/bin/sh
set -e

systemctl stop sqlite-demo.service
systemctl disable sqlite-demo.service
