# tests/sqlite-demo/Ubuntu_19.04/sqlite-demo/DEBIAN/postinst
#!/bin/sh
set -e

systemctl enable sqlite-demo.service
systemctl start sqlite-demo.service
