# tests/sqlite-demo/Ubuntu_19.04/sqlite-demo/DEBIAN/prerm
#!/bin/bash -e

systemctl stop sqlite-demo.service
systemctl disable sqlite-demo.service
